# install/templates/rc.local

/bin/sh /usr/local/share/mailserv/init.sh

