# install/templates/rc.local

/bin/sh /var/mailserv/scripts/mailserv_boot.sh
