# install/templates/rc.local

/bin/sh /usr/local/share/mailserver/init.sh

