# install/templates/rc.local
env LD_PRELOAD=/usr/lib/libpthread.so.18.0 /var/mailserv/scripts/mailserv_boot.sh
